# systemd/user/launch-taskbar-app@.service
[Unit]
Description=Launch taskbar app %i

[Service]
Type=oneshot
ExecStart=/home/lewis/Dev/config/bin/launch-taskbar-app.sh %i
